# supervisor/log_generator.ini
[program:log_generator]
command=/root/corgi/log_generator.py &
stdout_logfile=/root/corgi/output/log_generator.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true

stderr_logfile=/root/corgi/supervisor/log_generator-err.log
