# supervisor/push_bandwidth_filebeat.ini
[program:push_bandwidth_filebeat]
command=filebeat.sh -c /data/filebeat/push_bandwidth_filebeat.yml &
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true

stderr_logfile=/root/corgi/supervisor/push_bandwidth_filebeat-err.log
