# supervisor/seedpool_transport.ini
[program:seedpool_transport]
command=/root/corgi/seedpool_transport.py &
stdout_logfile=/root/corgi/output/strategy_executor.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true

stderr_logfile=/root/corgi/supervisor/seedpool_transport-err.log
