# supervisor/simulator.ini
[program:simulator]
command=/root/corgi/simulator.py &
stdout_logfile=/root/corgi/output/simulator.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true

stderr_logfile=/root/corgi/supervisor/simulator-err.log
