# supervisor/strategy_executor.ini
[program:strategy_executor]
command=/root/corgi/strategy_executor.py &
stdout_logfile=/root/corgi/output/strategy_executor.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true

stderr_logfile=/root/corgi/supervisor/strategy_executor-err.log